# server-only/minecraft.cron
# only for ramdisk
#*/5 * * * * /etc/init.d/minecraft.init disksaverun >/dev/null 2>/dev/null
20 */6 * * * /etc/init.d/minecraft backup >/dev/null 2>/dev/null
